# deploy/cc-sidecar.service
[Unit]
Description=CC Sidecar — Claude Code session watcher
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=mike
EnvironmentFile=-/etc/cc-sidecar/env
ExecStart=/home/mike/cc-sidecar/bin/cc-sidecar --config /etc/cc-sidecar/config.yaml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
