# deploy/cc-sidecar.service
[Unit]
Description=CC Sidecar — Claude Code session watcher
After=network.target nats-server.service
Wants=nats-server.service

[Service]
Type=simple
User=mike
ExecStart=/usr/local/bin/cc-sidecar
Restart=on-failure
RestartSec=5
Environment=CC_SIDECAR_NATS_URL=nats://localhost:4222

[Install]
WantedBy=multi-user.target
